# supervisor.conf
[program:ombra-backend]
command=/home/azureuser/Ombra/venv/bin/python -m uvicorn server:app --host 0.0.0.0 --port 8001
directory=/home/azureuser/Ombra/backend
user=azureuser
autostart=true
autorestart=true
stderr_logfile=/var/log/ombra-backend.err.log
stdout_logfile=/var/log/ombra-backend.out.log
environment=PATH="/home/azureuser/Ombra/venv/bin:%(ENV_PATH)s"
